# systemd.service.example
[Unit]
Description=Marcus Lion Chess Player Analyser
After=network.target

[Service]
User=YOUR_LINUX_USER
WorkingDirectory=/home/YOUR_LINUX_USER/chess-player-analyser
Environment="PATH=/home/YOUR_LINUX_USER/chess-player-analyser/.venv/bin"
ExecStart=/home/YOUR_LINUX_USER/chess-player-analyser/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
